# components/certbot/certbot-renew.timer
[Unit]
Description=Certbot renewal timer (twice daily)

[Timer]
OnCalendar=*-*-* 02,14:00:00
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
